# crontabs/www-data-example
# do daily/weekly/monthly maintenance
# min   hour    day     month   weekday command
* * * * * echo "I'm fine." >> /www/php82_cron_health_check.txt

# 将此文件重命名为希望执行此配置的用户，例如 www-data

# 每次重启容器后到 PHP 容器执行命令启动服务
# /usr/sbin/crond -b
